# deluge/rootfs/etc/services.d/deluge-web/run
#!/usr/bin/with-contenv bash

DELUGE_LOGLEVEL=${DELUGE_LOGLEVEL:-warning}

exec \
    /usr/bin/deluge-web -c /config/deluge\
	-d --loglevel=${DELUGE_LOGLEVEL}
